# fan-the-flames.service
# /etc/systemd/system/fan-the-flames.service

[Unit]
Description=Fan the Flames telnet fire splash screen
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/srv/fan-the-flames
ExecStart=/srv/fan-the-flames/venv/bin/python /srv/fan-the-flames/telnet_server.py --host 0.0.0.0 --port 23
Restart=always
RestartSec=5

# Run as a throwaway unprivileged user instead of root...
DynamicUser=yes
# ...but still let it bind the privileged telnet port (23)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=default.target
