# fan-the-flames.service
# /etc/systemd/system/fan-the-flames.service

[Unit]
Description=Fan the Flames telnet fire splash screen

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/opt/fan-the-flames/telnet_server.py --host 0.0.0.0 --port 23

[Install]
WantedBy=default.target
